# src/Monik.Service/Monik.Service.service
[Unit]
Description=Monik.Service

[Service]
Type=notify
ExecStart=/usr/sbin/Monik.Service

[Install]
WantedBy=multi-user.target